# mntbooks.service
[Unit]
Description=MNT Bookkeeping

[Service]
ExecStart=/home/mntbooks/mntbooks/start.sh

[Install]
WantedBy=multi-user.target
